# conf/LBCAlerte.cron
*/5 * * * * www-data /usr/bin/php -f #DESTDIR#/check.php